# scripts/com.daynightcycle.updater.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.daynightcycle.updater</string>

    <key>ProgramArguments</key>
    <array>
        <string>PYTHON_PATH_PLACEHOLDER</string>
        <string>PROJECT_DIR_PLACEHOLDER/scripts/update_schedule.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>PROJECT_DIR_PLACEHOLDER</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>PROJECT_DIR_PLACEHOLDER/logs/updater.log</string>

    <key>StandardErrorPath</key>
    <string>PROJECT_DIR_PLACEHOLDER/logs/updater.error.log</string>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
